# deploy/deb-server/DEBIAN/postinst
#!/bin/bash
set -e

echo "Atolye Platform Server yapılandırılıyor..."

# Bağımlılıkları kontrol et ve yerel mimariye göre derle/rebuild yap
if command -v npm &> /dev/null; then
    echo "better-sqlite3 modülü yerel mimari için yapılandırılıyor..."
    cd /opt/atolye-server/server
    # Eşleşmeyen native binary'leri temizle ve yerel olarak yeniden derlemeyi dene
    # Hata durumunda kurulumun yarıda kalmaması için exit code 0 ile devam edilmesini sağlıyoruz
    npm rebuild better-sqlite3 || echo "⚠️ better-sqlite3 yerel olarak derlenemedi. Sunucu hibrit JSON modunda çalışacaktır."
fi

# Servis dosyasını etkinleştir
systemctl daemon-reload
systemctl enable atolye-server.service
systemctl restart atolye-server.service

echo "✅ Atolye Platform Server başarıyla yüklendi ve başlatıldı."
echo "🌐 Ağda 'atolye-server.local' olarak bulunabilir."
